# conf/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn.sock
User=gunicorn

[Install]
WantedBy=sockets.target
